# build/apps/core-api/supervisor/conf.d/kafka-search.conf
[program:kafka-search-events]
directory=/var/www
command=php artisan kafka:listen --topic=searchEvents
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/kafka-search-events.log
stderr_logfile=/var/log/supervisor/kafka-search-events.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=9
environment=HOME="/var/www"
