# build/apps/core-api/supervisor/conf.d/kafka-booking.conf
[program:kafka-bookings]
directory=/var/www
command=php artisan kafka:listen --topic=bookings
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/kafka-bookings.log
stderr_logfile=/var/log/supervisor/kafka-bookings.err
priority=10
environment=HOME="/var/www"
