# build/apps/core-api/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm.err

